# platform/pi/services/mavlink-router.service
[Unit]
Description=Mavlink Router
Wants=network.target
After=syslog.target network-online.target

[Service]
Type=exec
ExecStart=mavlink-routerd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
